# isos/appliance/vic.target
[Unit]
Description=vSphere Integrated Containers
Requires=network-online.target
After=network-online.target
AllowIsolate=yes